# files/first-time-checker/first-time-checker.service
[Unit]
Description=First-time system setup checker
After=graphical-session.target

[Service]
Type=oneshot
ExecStart=/usr/bin/first-time-checker.sh
RemainAfterExit=no

[Install]
WantedBy=default.target
